# SEM06/resurse/systemd/backup.service
# Systemd Service pentru Backup (folosit cu backup.timer)
# Salvează ca: /etc/systemd/system/backup.service

[Unit]
Description=Automated Backup Service
Documentation=man:backup(1)
After=network.target local-fs.target

[Service]
Type=oneshot

# User și group (ajustează pentru setup-ul tău)
User=root
Group=root

# Working directory
WorkingDirectory=/opt/capstone/backup

# Script principal cu argumente
ExecStart=/opt/capstone/backup/backup.sh create \
    --source=/var/www \
    --dest=/backup \
    --type=incremental \
    --compress=gz \
    --retention=7

# Pre și post scripts
ExecStartPre=/bin/mkdir -p /backup
ExecStartPost=/usr/bin/logger "Backup completed successfully"

# Timeout (backup-urile pot dura)
TimeoutStartSec=3600

# Logging
StandardOutput=append:/var/log/backup/backup.log
StandardError=append:/var/log/backup/backup_error.log

# Environment
Environment="BACKUP_NOTIFY_EMAIL=admin@example.com"
Environment="LOG_LEVEL=INFO"

# Security (mai relaxat pentru backup care citește tot sistemul)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
